# etc/wii-u-gc-adapter.service
[Unit]
Description=Wii U GC Adapter

[Service]
TimeoutStartSec=infinity
ExecStart=sh -c "${HOME}/.local/bin/wii-u-gc-adapter"
